# etc/beta/supervisord/aspnet-main-web.ini
[program:aspnet-main-web]
user=web
directory=/home/web/aspnet-main/src/App/bin/Release/netcoreapp3.0/linux-x64/publish
command=/bin/bash -c ". /home/web/env/aspnet-main-web-beta.env && ./App"
environment=ASPNETCORE_ENVIRONMENT="Beta",ASPNETCORE_URLS="http://localhost:9500"
redirect_stderr=true
stdout_logfile=/home/web/logs/aspnet-main-web.log
stopasgroup=true
killasgroup=true
